# launchd/com.goodproject.source-scout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.goodproject.source-scout</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/reify/Classified/goodcapital_landing/launchd/run_source_scout.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/reify/Classified/goodcapital_landing</string>

  <key>StartInterval</key>
  <integer>14400</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/reify/Classified/goodcapital_landing/diagnostics/source-scout-launchd.out</string>

  <key>StandardErrorPath</key>
  <string>/Users/reify/Classified/goodcapital_landing/diagnostics/source-scout-launchd.err</string>
</dict>
</plist>
